# cron/.config/.crontab
* * * * * /usr/bin/offlineimap -o -u syslog
5 * * * * $HOME/.config/vdirsyncer/.venv/bin/vdirsyncer sync
