# cron/.config/.crontab
*/15 * * * * rsync ~/src/ truenas.localdomain:src/
