# cron/.config/.crontab
* * * * * /usr/bin/offlineimap -o -u syslog
5 * * * * /usr/bin/vdirsyncer sync
