# cron/.config/.crontab
* * * * * /usr/bin/offlineimap -o -u syslog
