# systemd/.config/systemd/user/offlineimap.timer
[Unit]
Description=Offlineimap Query Timer
Requires=offlineimap.service

[Timer]
OnUnitInactiveSec=30s
Unit=offlineimap.service

[Install]
WantedBy=mail.target
